# 99-razer-blackwidow-ultimate.rules
ACTION=="add", SUBSYSTEM=="input", ATTRS{idVendor}=="1532", ATTRS{idProduct}=="010d", RUN+="/usr/bin/razer-blackwidow-ultimate.py"
